# home/dot_config/systemd/user/chezmoi-update.service.d/override.conf
[Unit]
ConditionPathExists=!%h/chezmoi-update-inhibit
