# home/dot_config/systemd/user/chezmoi-diff.timer
[Unit]
Description=Periodically count chezmoi differences

[Timer]
OnBootSec=5m
OnUnitInactiveSec=5m

[Install]
WantedBy=timers.target
